# roles/stream/files/cron.daily/certbot.cron
#!/bin/sh
/opt/certbot-auto renew --quiet
