# roles/stream/files/cron.daily/letsencrypt
#!/bin/sh
letsencrypt renew --quiet --renew-hook /opt/gen_icecast_cert.sh
