# roles/certbot/files/letsencrypt-renewer.timer
[Unit]
Description=Renew Let's Encrypt certificates

[Timer]
OnUnitActiveSec=12h

[Install]
WantedBy=multi-user.target
