# roles/certbot/files/letsencrypt-renewer.timer
[Unit]
Description=Renew Let's Encrypt certificates

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=multi-user.target
