# files/media-local\x2dstorage.mount
[Mount]
What=/dev/storage/local-storage
Where=/media/local-storage
Type=ext4
